# final_project/main.cron
0 7 * * * python main.py > logs.log 2>&1
